# provisioning/roles/kubera/templates/etc/profile.d/kubera.sh.j2
export CORS_ORIGINS="{{hostvars['kubera'].private_ipv4}}"
